# infra/systemd/infinda-api.service
[Unit]
Description=INFINDA Django API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=infinda
Group=infinda
WorkingDirectory=/opt/infinda/apps/api
EnvironmentFile=/opt/infinda/infra/env/infinda.app.env
ExecStart=/opt/infinda/apps/api/.venv/bin/python manage.py runserver 127.0.0.1:8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
